# _superninja_startup.conf
autostart=true
[program:quizforge_server]
command=python3 /workspace/quizforge/server.py
directory=/workspace/quizforge
autostart=true
autorestart=true
startretries=10
startsecs=3
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/workspace/quizforge/server.log
stdout_logfile_maxbytes=5MB
environment=PORT=8085
